# deploy/systemd/rdp-casdoor.service
[Unit]
Description=RDP Casdoor Service
Documentation=https://github.com/casdoor/casdoor
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=rdp
Group=rdp
WorkingDirectory=/opt/casdoor
ExecStart=/opt/casdoor/casdoor
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5s

LimitNOFILE=65536

StandardOutput=journal
StandardError=journal
SyslogIdentifier=rdp-casdoor

[Install]
WantedBy=multi-user.target
