# deploy/systemd/rdp-casdoor.service
[Unit]
Description=RDP Casdoor Authentication Service
Documentation=https://github.com/casdoor/casdoor
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=rdp
Group=rdp
WorkingDirectory=/opt/casdoor
ExecStart=/opt/casdoor/casdoor
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5s
StartLimitInterval=60s
StartLimitBurst=3

# Environment
EnvironmentFile=-/opt/casdoor/conf/app.conf

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/casdoor/logs /opt/casdoor/files

# Resource limits
LimitNOFILE=65536
LimitNPROC=65536

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rdp-casdoor

[Install]
WantedBy=multi-user.target rdp.target
